# .config/systemd/user/novnc.service
[Unit]
Description=noVNC proxy

[Service]
ExecStart=/usr/bin/novnc_proxy
Type=simple

[Install]
WantedBy=default.target
